# deploy/galacticmap-autoupdate.service
[Unit]
Description=Galactic Map self-update (pull from GitHub + rebuild)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
# Optional: holds GH_TOKEN (and any GM_* overrides). The leading - means "skip if absent" — a public repo
# needs no token, so the file is optional.
EnvironmentFile=-/etc/galactic-map-autoupdate.env
ExecStart=/usr/local/bin/galacticmap-autoupdate
TimeoutStartSec=1200
